# rootfs/etc/inittab
::sysinit:/etc/init.d/rc
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::restart:/sbin/init
tty1::respawn:/bin/getty 115200 tty1 linux
tty2::respawn:/bin/getty 115200 tty2 linux
tty3::respawn:/bin/getty 115200 tty3 linux
tty4::respawn:/bin/getty 115200 tty4 linux
tty5::respawn:/bin/getty 115200 tty5 linux
tty6::respawn:/bin/getty 115200 tty6 linux
ttyS0::respawn:/bin/getty 115200 ttyS0
